# ansible/roles/vpsbridges/templates/dehydrated.timer
[Unit]
Description=Run dehydrated certificate refresh

[Timer]
OnCalendar=Mon 13:00

[Install]
WantedBy=timers.target

