# ansible/roles/wiki/templates/node-wiki.service.j2
[Service]
ExecStart=/usr/bin/node wiki start
Restart=always
StandardOutput=syslog
StandardError=syslog
User={{ wiki_user }}
Group={{ wiki_user }}
RestartSec=15s

[Install]
WantedBy=multi-user.target
