# ansible/roles/ooni-orchestrate/templates/ooni-orchestrate.service
[Unit]
Description=Proteus Orchestrate service

[Service]
User={{ orchestra_user }}
Group={{ orchestra_group }}
ExecStart={{ orchestra_path }}/orchestra-orchestrate --config {{ orchestra_config_path }}/orchestra-orchestrate.toml start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
