# ansible/roles/monitoring/files/jupycron.timer
[Unit]
Description=Run jupycron every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
