# ansible/roles/letsencrypt-legacy/templates/monitor_certs.service
[Unit]
Description=Run monitor_certs.py, activated by monitor_certs.timer

[Service]
Type=oneshot
ExecStart=/etc/ooni/monitor_certs.py

