# ansible/roles/node_exporter/templates/smartmon_metrics.timer
[Unit]
Description=smartmon metrics

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
