# ooni-sysadmin-ng/ooni-probe/_ooniprobe-bridge-run.timer
[Unit]
Description=Daily ooniprobe-bridge-reachability

[Timer]
OnBootSec=10min
OnUnitActiveSec=1day

[Install]
WantedBy=multi-user.target
